# workshops/partner-intro/workshop/profile.d/01-workshop-env.sh
#!/bin/bash
# Symlink exercises into home directory for easy access from terminal
ln -sfn /opt/workshop/content/exercises ~/exercises
